# src/monitor/services/systemd/tedge-container-monitor.service
[Unit]
Description=thin-edge.io container status monitor
After=mosquitto.service

[Service]
ExecStart=/usr/bin/tedge-container-monitor
Restart=always
RestartSec=30
EnvironmentFile=-/etc/tedge-container-plugin/env
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
